# application/rootfs/etc/inittab
#
# inittab       This file describes how the INIT process should set up
#               the system during the initial booting process.

# System initialization starts here


::sysinit:/etc/rc.d/rc.sysinit
::respawn:/sbin/getty -L 115200 ttyAMA0 vt102
::shutdown:/etc/rc.d/rc.shutdown
::ctrlaltdel:/bin/false
::respawn:/usr/sbin/syslog-ng -f /etc/syslog-ng/syslog-ng.conf -p /var/run/syslogd.pid -F -R /opt/pme/pram/syslog-ng.persist > /dev/null 2>&1
::respawn:/opt/pme/apps/fsync
